# 5.5/rootfs/etc/services.d/solr/run
#!/bin/sh

exec s6-setuidgid wodby /opt/solr/bin/solr start -f
