# 5.5/rootfs/etc/services.d/solr/run
#!/bin/sh

exec /opt/solr/bin/solr start
